# supervisor.conf
[program:traefik]
command=traefik --configFile=/etc/traefik/traefik.toml
autostart=true
autorestart=true
startsecs=3